# btrfs-scrub@.service.in
[Unit]
Description=Scrub a BTRFS filesystem.
Documentation=man:btrfs-scrub(8)
Conflicts=btrfs-balance@.service

[Service]
ExecStart=/bin/btrfs scrub start -B -d  -c idle %I
ExecStop=/bin/btrfs scrub cancel %I
ExecStopPost=/bin/btrfs scrub status %I

# systemd.exec settings #
IOSchedulingClass=idle
CPUSchedulingPolicy=idle
StandardOutput=journal
StandardError=journal
SyslogIdentifier=BTRFS-Scrub


