# btrfs-scrub@.timer.in
[Unit]
Description=Timer for scrubbing a btrfs filesystem
Documentation=man:btrfs-scrub(8)

[Timer]
OnCalendar=monthly
Persistent=true