# btrfs-balance@.service.in
[Unit]
Description=Balance a BTRFS filesystem.
Documentation=man:btrfs-balance(8)
Conflicts=btrfs-scrub@.service

[Service]
ExecStart=/bin/btrfs balance start -v %I
ExecStop=/bin/btrfs balance cancel %I
ExecStopPost=/bin/btrfs balance status -v %I

# systemd.exec settings #
IOSchedulingClass=idle
CPUSchedulingPolicy=idle
StandardOutput=journal
StandardError=journal
SyslogIdentifier=BTRFS-Balance


